# conf/rc.local
#!/bin/sh -e
service postfix start
exit 0
